# core/src/main/packages/macos/scripts/banalytics-box.service
[Unit]
Description=BanalyticsBox

[Service]
ExecStart=/opt/banalytics/banalytics-box-start.sh
User=banalytics
Group=banalytics

[Install]
WantedBy=multi-user.target