# core/src/main/packages/ubuntu/pulseaudio.service
[Unit]
Description=PulseAudio system server
# DO NOT ADD ConditionUser=!root

Type=notify
Exec=pulseaudio --start --log-target=journal
Restart=on-failure

[Install]
WantedBy=multi-user.target