# core/src/main/packages/old-linux-backup/scripts/banalytics-box.service
[Unit]
Description=BanalyticsBox

[Service]
ExecStart=/opt/banalytics/banalytics-box-start.sh
User=banalytics-service
Group=banalytics-service

[Install]
WantedBy=multi-user.target